# proxy/proxy.conf
[supervisord]
nodaemon=true

[program:proxy]
command=/app/proxy
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/proxy.err.log
stdout_logfile=/var/log/supervisor/proxy.out.log
