# dev/voideditor/void_settings.conf
[supervisord]
nodaemon=true

[program:void_settings]
command=bash -c "mkdir -p ~/.config/Void/User && cp /app/Void_User_settings.json ~/.config/Void/User/settings.json"
autostart=true
autorestart=false
user=%(ENV_USER)s
directory=/home/%(ENV_USER)s
stderr_logfile=/var/log/void_settings.err.log
stdout_logfile=/var/log/void_settings.out.log